# printflow.service
[Unit]
Description=PrintFlow Automation Server
After=network.target

[Service]
Type=simple
User=fotos
WorkingDirectory=/home/fotos/printflow
Environment=BIND_HOST=0.0.0.0
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=append:/home/fotos/printflow/server.log
StandardError=append:/home/fotos/printflow/server.log

[Install]
WantedBy=multi-user.target
